# systemd/voice-bridge.service
[Unit]
Description=Voice Bridge — Jabra speakerphone to gateway voice client
# This is a *user* unit — installed under ~/.config/systemd/user and
# managed with `systemctl --user`. It runs as the invoking user (no
# `User=`) and inherits that user's group memberships, so the user
# must already be in `audio` for ALSA and `plugdev` for /dev/hidraw*
# (granted by the matching udev rule). Enable linger
# (`sudo loginctl enable-linger <user>`) so the unit starts at boot
# without a login session.
After=sound.target network-online.target
Wants=network-online.target
# Back off if the bridge crash-loops fast (StartLimit* live in [Unit]).
StartLimitIntervalSec=60
StartLimitBurst=5

[Service]
Type=simple
# %h is systemd's specifier for the running user's home dir (user unit).
# Edit the Services/voice-bridge suffix if your checkout lives elsewhere.
WorkingDirectory=%h/Services/voice-bridge
# Diagnostic: surface aplay's stderr (underrun!!! lines) to the journal.
Environment=VOICE_BRIDGE_DEBUG_AUDIO=1
ExecStart=%h/Services/voice-bridge/.venv/bin/python %h/Services/voice-bridge/voice-bridge.py
Restart=on-failure
RestartSec=5
# The bridge handles SIGTERM/SIGINT cleanly: finishes the in-flight
# turn, stops the HID thread, exits. The HID monitor's reconnect loop
# parks on a threading.Event between attempts, so a missing Jabra
# costs ~0 CPU — no need for the unit to stop and restart on unplug.
KillSignal=SIGTERM
TimeoutStopSec=10

[Install]
# default.target is the user-mode equivalent of multi-user.target.
WantedBy=default.target
